# v1/docker/etc/services.d/postgresql/run
#!/usr/bin/execlineb -P

define PGDATA /var/lib/postgresql
export PGDATA ${PGDATA}

define POSTGRES_USER postgres

s6-setuidgid ${POSTGRES_USER} postgres
